# fluent-bit/files/service.upstart.conf.j2
{% from "fluent-bit/map.jinja" import fluent_bit with context -%}
description "Fluentbit log collection service"

start on started network-services
stop on runlevel [!2345]
setgid nogroup

# Set reasonable file limits
limit nofile 65536 65536

# Restart automatically on crash
respawn

console log
{%- if fluent_bit.env is defined %}
# These environment variables are generated by salt.
{%- for key, value in fluent_bit.env.items() %}
env {{ key }}={{ value }}
{%- endfor %}
{%- endif %}

exec /usr/local/bin/fluent-bit -c /etc/fluent-bit/fluent-bit.conf -R /etc/fluent-bit/parsers.conf
